# docs/systemd/payment-gateway-worker.service.example
[Unit]
Description=Payment Gateway Queue Worker
After=network.target

[Service]
User=www-data
Group=www-data
Restart=always
RestartSec=5
WorkingDirectory=/path/to/project
ExecStart=/usr/bin/php /path/to/project/artisan queue:work --queue=webhooks,payouts,default --tries=3 --timeout=120 --memory=256
StandardOutput=append:/path/to/project/storage/logs/worker.log
StandardError=append:/path/to/project/storage/logs/worker-error.log

[Install]
WantedBy=multi-user.target

